# service/rsyslog/run
#!/bin/sh
set -e

# Run syslogd
exec rsyslogd -n
