# service/postfix/run
#!/bin/sh
set -e

# Run postfix master
/usr/libexec/postfix/master -c /etc/postfix -d 2>&1
